# deploy/question-bank.service.template
[Unit]
Description=<SERVICE_NAME>
After=network.target

[Service]
WorkingDirectory=<APP_ROOT>
Environment=QUESTION_BANK_DB_PATH=<SERVER_DB_PATH>
Environment=QUESTION_BANK_ASSETS_DIR=<SERVER_ASSETS_DIR>
Environment=QUESTION_BANK_RAW_DIR=<SERVER_RAW_DIR>
ExecStart=<APP_ROOT>/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port <SERVICE_PORT>
Restart=always
RestartSec=5
User=<SERVICE_USER>
Group=<SERVICE_GROUP>

[Install]
WantedBy=multi-user.target
